# system_files/usr/lib/systemd/system/bootc-fetch-apply-updates.timer.d/updates.conf
[Timer]
# Clear previous timers
RandomizedDelaySec= OnBootSec= OnBootSec=10m OnUnitInactiveSec=10m
