# system_files/usr/lib/systemd/system/bootc-fetch.timer
[Unit]
Description=Fetch bootc image ASAP
[Timer]
# OnCalendar=*-*-* 12:30:00
OnUnitActiveSec=10m
Persistent=true

[Install]
WantedBy=timers.target